# board/f+s/fsimx6/inittab
# /etc/inittab
#
# Note: BusyBox init doesn't support runlevels.  The runlevels field is
# completely ignored by BusyBox init. If you want runlevels, use
# sysvinit.
#
# Format for each entry: <id>:<runlevels>:<action>:<process>
#
# id        == tty to run on, or empty for /dev/console
# runlevels == ignored
# action    == one of sysinit, respawn, askfirst, wait, and once
# process   == program to run

# Startup the system
::sysinit:/bin/mount -t proc proc /proc
#::sysinit:/bin/mount -t tmpfs -o size=64k,mode=0755 tmpfs /dev
#::sysinit:/bin/mount -o remount,rw / # REMOUNT_ROOTFS_RW
::sysinit:/bin/mkdir -p /dev/pts
::sysinit:/bin/mkdir -p /dev/shm
::sysinit:/bin/echo /sbin/mdev > /proc/sys/kernel/hotplug
::sysinit:/bin/mount -a
::sysinit:/bin/hostname -F /etc/hostname
# now run any rc scripts
::sysinit:/etc/init.d/rcS

# Put a getty on the login_tty. /etc/init.d/S01fssetup creates this as a
# symbolic link to the device given in /sys/bdinfo/login_tty, usually the
# serial port; we also set the speed that is given in /sys/bdinfo/login_speed
ttymxc1::respawn:/sbin/getty -L ttymxc1 115200 vt100
#login_tty::respawn:/sbin/getty -L login_tty $(/bin/cat /sys/bdinfo/login_speed) vt100

# Stuff to do for the 3-finger salute
::ctrlaltdel:/sbin/reboot

# Stuff to do before rebooting
null::shutdown:/etc/init.d/rcK
null::shutdown:/bin/umount -a -r
null::shutdown:/sbin/swapoff -a
