# board/f+s/common/etc/profile.d/fsalias.sh
#!/bin/sh
# F&S additional and useful aliases

alias cp="cp -i"
alias mv="mv -i"
alias rm="rm -i"
alias Rm="rm -f"
alias l="ls -lF"
alias ll="ls -lF"
alias ls="ls -F"
alias ..="cd .."
alias pd='pushd'
