# board/f+s/common/etc/profile.d/weston.sh
#!/bin/sh
# Export environment variable for weston

export XDG_RUNTIME_DIR=/var/run/weston
