# board/f+s/common/etc/profile.d/xserver.sh
#!/bin/sh
# XServer setup to be run from /etc/profile.d

export DISPLAY=:0
export XAUTHORITY=/.Xauthority
